# SPN/ENTER/etc/profile.d/conda.csh
setenv CONDA_EXE "/home/bamichel/Documents/MILP_cryptanalysis/SPN/ENTER/bin/conda";
setenv _CONDA_EXE "/home/bamichel/Documents/MILP_cryptanalysis/SPN/ENTER/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/home/bamichel/Documents/MILP_cryptanalysis/SPN/ENTER/bin/python";
setenv _CONDA_ROOT "/home/bamichel/Documents/MILP_cryptanalysis/SPN/ENTER";

source "/home/bamichel/Documents/MILP_cryptanalysis/SPN/ENTER/lib/python3.12/site-packages/conda/shell/etc/profile.d/conda.csh"